# api-service/supervisor.conf
[program:docker-service]
command=docker build -t flask-api .
command=docker run -p 8000:8000 flask_api:latest
directory=/home/api/api-source
user=user
